# root/etc/systemd/system/dedalo.service.d/disable_eth_offload.conf
[Service]
ExecStartPre=-/usr/libexec/nethserver/dedalo-disable-ethernet-offload
